# src/test/resources/META-INF/test-6.cron
* 47 0 * * * /etc/acme/scripts/cron-run.sh >> /var/log/acme/cron.log 2>&1;
