# src/main/resources/META-INF/hello-world.cron
0/5 * * * * ? /usr/local/bin/hello-world.sh
